# rootfilesystem/etc/supervisor/service.d/swoole.conf
[supervisord]
user = root

[program:swoole]
command = php /var/www/server.php
user = root
autostart = true
autorestart = true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/1
stderr_logfile_maxbytes=0
